# overlord-bridge.service
[Unit]
Description=[Ben] Overlord Telegram bridge (drive Claude Code from your phone)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/ben/Projects/overlord-bridge
ExecStart=/home/ben/Projects/overlord-bridge/.venv/bin/python /home/ben/Projects/overlord-bridge/bridge.py
Restart=on-failure
RestartSec=5
Environment=PYTHONUNBUFFERED=1
# The Agent SDK shells out to the `claude` CLI, and the Codex brain shells out
# to `codex app-server`; both live in npm-global, so keep venv + npm-global on
# PATH. The brain is chosen via OVERLORD_BRAIN in .env — no unit change needed
# to switch (just restart the service). codex uses ~/.codex/auth.json, no extra
# env required here.
Environment="PATH=/home/ben/Projects/overlord-bridge/.venv/bin:/home/ben/.npm-global/bin:/usr/local/bin:/usr/bin:/bin"

[Install]
WantedBy=default.target
